# docker/nuxt/s6/etc/services.d/elasticsearch/run
#!/usr/bin/with-contenv sh

cd /opt/nuxt/elasticsearch
exec node server.js
