# docker/nuxt/s6/etc/s6-overlay/s6-rc.d/nuxt3/run
#!/command/with-contenv sh

cd /opt/nuxt3
if [ ! -f "/opt/nuxt3/dist/server/index.mjs" ]; then
    pnpm build
fi
exec node dist/server/index.mjs
