# docker/nuxt/s6/etc/s6-overlay/s6-rc.d/elasticsearch/run
#!/command/with-contenv sh

cd /opt/nuxt/elasticsearch
exec node server.mjs
