# docker/nuxt/s6/etc/services.d/nuxt/run
#!/usr/bin/with-contenv sh

cd /opt/nuxt
if [ ! -d "/opt/nuxt/.nuxt/dist/" ]; then
    nuxt build
fi
exec nuxt start
